# ssh1d.socket
[Socket]
ListenStream = 22
BindIPv6Only = both

[Install]
WantedBy = sockets.target
