# ssh1d.service
[Unit]
Description = ssh1d - fake SSH-1 daemon

[Service]
ExecStart = /usr/bin/ssh1d
